# data/services/bootstrap/watch-iri-tls-certs.path
[Unit]
Description=Watch for IRI TLS certificate files
Documentation=man:systemd.path(5)

[Path]
PathExists=/opt/openshift/tls/internal-release-image.crt
PathExists=/opt/openshift/tls/internal-release-image.key
Unit=reconfigure-local-registry-iri-tls.service

[Install]
WantedBy=multi-user.target
