# data/services/common/start-local-registry.service
[Unit]
Description=Local Registry
Wants=network.target

[Service]
ExecStart=/usr/local/bin/start-local-registry.sh
Restart=on-failure

[Install]
WantedBy=multi-user.target
